# homelab/remote-desktop/vnc/x11vnc.service
# systemd USER service: start x11vnc (X11 desktop) bound to the tailnet.
# Install:
#   mkdir -p ~/.config/systemd/user
#   cp x11vnc.service ~/.config/systemd/user/
#   # edit ExecStart path below to where you cloned this repo
#   systemctl --user daemon-reload
#   systemctl --user enable --now x11vnc.service
# Inspect:
#   systemctl --user status x11vnc.service
#   journalctl --user -u x11vnc.service -f
#
# Prereq: create the password file first:
#   x11vnc -storepasswd ~/.config/homelab/vncpasswd
#
# Note: this is a USER unit on purpose — it needs your X session's $DISPLAY.
# If your display isn't :0, set X11_DISPLAY in the [Service] Environment below.

[Unit]
Description=x11vnc VNC server (bound to Tailscale only)
After=graphical-session.target tailscaled.service
PartOf=graphical-session.target

[Service]
Type=simple
# EDIT this path to the actual location of the script on your PC:
ExecStart=%h/open/homelab/remote-desktop/vnc/start-vnc.sh
Environment=VNC_SERVER=x11vnc
# Uncomment / adjust if your X display is not :0
# Environment=X11_DISPLAY=:0
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical-session.target
